# turbo.service
[Unit]
Description=Turbo server

[Service]
ExecStart=/home/turbo/turbo
User=turbo
Group=turbo
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
